# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /usr/src/app/uwsgi.ini --chdir /usr/src/app/src

[program:nginx-app]
command = /usr/sbin/nginx
